# systemd/disk-sentinel.service
[Unit]
Description=Daily disk sentinel check

[Service]
Type=oneshot
ExecStart=%h/.local/libexec/noba/disk-sentinel.sh

[Install]
WantedBy=multi-user.target
