# systemd/log-rotator.service
[Unit]
Description=Weekly log rotation

[Service]
Type=oneshot
ExecStart=%h/.local/libexec/noba/log-rotator.sh

[Install]
WantedBy=multi-user.target
