# backup-20260316-084907/systemd/disk-sentinel.timer
[Unit]
Description=Daily disk sentinel check

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
