# backup-20260316-084907/systemd/system-report.timer
[Unit]
Description=Weekly system report

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
