# backup-20260316-082744/systemd/log-rotator.timer
[Unit]
Description=Weekly log rotation

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
